# shellscript.target
[Unit]
Description=mutant docker
Wants=shellscript.timer


[Install]
WantedBy=multi-user.target
