# shellscript.service
[Unit]
Description=My Shell Script

[Service]
ExecStart=/home/vagrant/repo3/pull_artic_dockers.sh


[Install]
WantedBy=multi-user.target
