# shellscript.service
[Unit]
Description=My Shell Script

[Service]
ExecStart=/home/vagrant/repo3/test.sh


[Install]
WantedBy=multi-user.target
